# meeting-electron/packages/meeting-kit-electron/src/com.netease.nmc.rooms.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.netease.nmc.rooms.agent</string>
  <key>ProgramArguments</key>
  <array>
      <string>/usr/bin/open</string>
      <string>-a</string>
      <string>/Applications/网易会议 Rooms.app</string>
  </array>
  <key>KeepAlive</key><true/>
</dict>
</plist>